# etc/supervisor/conf.d/crond.conf
[program:crond]
command=/usr/sbin/crond -f -l 2
user=root
autostart=true
autorestart=true
priority=300
stdout_logfile=/var/log/crond.log
stderr_logfile=/var/log/crond.err

